# scripts/Examples/bellTest/bellTest.service
[Unit]
Description=Prueba del timer de systemd

[Service]
Type=oneshot
ExecStart=/home/leandro/Documents/Gugler/bellTest/bellTest.sh
